# assets/config/supervisor/7-extra.conf
[program:extra]
command=bash -c "${EXTRA_COMMAND:-sleep infinity}"
priority=7
autorestart=true
stderr_logfile=/proc/self/fd/1
