# assets/config/supervisor/7-grass.conf
[program:Grass]
command=grass
priority=7
stderr_logfile=/proc/self/fd/2
